# services/platform-core-detection.service.in
[Unit]
Description=@PROJECT_NAME@ - core detection service
DefaultDependencies=no
Before=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/pr-detect --step=core --output-file=@DETECTION_OUTPUT_DIR@/core --json-file=@DETECTION_OUTPUT_DIR@/core.json

[Install]
WantedBy=systemd-modules-load.service

