# services/platform-devices-customize.service.in
[Unit]
Description=@PROJECT_NAME@ - devices customization service
DefaultDependencies=no
After=platform-devices-detection.service
Requires=platform-devices-detection.service

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/pr-customize --step=devices 

[Install]
WantedBy=basic.target

